# systemd/voice-to-text-deepgram.service
[Unit]
Description=Voice-to-Text (Push-to-Talk, Deepgram)
After=graphical-session.target

[Service]
Type=simple
WorkingDirectory=@REPO_DIR@
ExecStart=@REPO_DIR@/.venv/bin/python @REPO_DIR@/start_push_to_talk_deepgram.py
Restart=on-failure
RestartSec=2

# Required for X11-based global key listening and text insertion.
Environment=DISPLAY=:0
Environment=XAUTHORITY=%h/.Xauthority

[Install]
WantedBy=default.target